# pkgs/knot/root/etc/s6-overlay/s6-rc.d/knot/run
#!/command/with-contenv sh

exec s6-setuidgid git /bin/knot server
